# apilib/src/multicom/src/embx/tests/transport/harness/ramfs/etc/inittab
::sysinit:/etc/rc.sysinit
::askfirst:-/bin/sh
